# ops/systemd/pulsar-backup.timer
[Unit]
Description=Daily Pulsar SQLite backup

[Timer]
OnCalendar=*-*-* 03:15:00
RandomizedDelaySec=15m
Persistent=true
Unit=pulsar-backup.service

[Install]
WantedBy=timers.target
