# ops/systemd/pulsar-update.service
[Unit]
Description=Update Pulsar 2.0 from GitHub with validation and rollback
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pulsar-update
TimeoutStartSec=infinity
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
