# ops/systemd/pulsar-web.service
[Unit]
Description=Pulsar 2.0 Next.js web process
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pulsar
Group=pulsar
WorkingDirectory=/opt/pulsar/current
EnvironmentFile=/etc/pulsar/pulsar.env
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3100
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryMax=900M
TasksMax=256

[Install]
WantedBy=multi-user.target
