# ops/systemd/remnawave-backup.service
[Unit]
Description=Create and verify a Remnawave PostgreSQL backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/remnawave-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
