# deploy/systemd/pulsar-backup.timer
[Unit]
Description=Back up Pulsar SQLite every six hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:15:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
