# deploy/systemd/pulsar-web.service
[Unit]
Description=Pulsar 2.0 web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pulsar
Group=pulsar
WorkingDirectory=/opt/pulsar/current
EnvironmentFile=/etc/pulsar/pulsar.env
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
ExecStart=/usr/bin/node --max-old-space-size=768 .next/standalone/server.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pulsar
MemoryHigh=900M
MemoryMax=1200M
TasksMax=128
UMask=0077
OOMPolicy=stop
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
